# profile/profile.d/30-npm.sh
NPM_PACKAGES="${HOME}/.local/lib/node-packages"
varmunge PATH "$NPM_PACKAGES/bin"
varmunge MANPATH "$NPM_PACKAGES/share/man"

unset NPM_PACKAGES

